# scripts/knowledge-synthesize@.timer.in
# systemd USER timer template — generated by scripts/install.sh, which fills in
# @SYNTH_ONCALENDAR@ from KNOWLEDGE_SYNTHESIZE_ONCALENDAR (default: weekly). This is the heavy
# whole-corpus pass — keep it rare. Change the cadence via that env var (or .env), then re-run
# install.sh.
[Unit]
Description=Scheduled vault synthesis pass

[Timer]
OnCalendar=@SYNTH_ONCALENDAR@
# Run on next boot/login if the machine was off when a tick was due.
Persistent=true
# Spread out so it doesn't fire at the same instant as other timers.
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
