# scripts/knowledge-compile.service.in
# systemd USER service template — generated into ~/.config/systemd/user/knowledge-compile.service
# by scripts/install.sh, which fills in @TOOLS_REPO@ (this repo, where the worker lives) and
# @VAULT_REPO@ (the vault repo it compiles). Edit this template, then re-run install.sh.
#
# Runs one ephemeral inbox->wiki compile. Started by two units that share this single
# service instance (so systemd serializes them — the cross-process lock):
#   - knowledge-compile.timer → scheduled run (cadence via KNOWLEDGE_COMPILE_ONCALENDAR)
#   - knowledge-compile.path  → on-demand run when the MCP drops inbox/.compile/request
#
[Unit]
Description=Compile the vault inbox into the wiki (ephemeral Claude run)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Point the worker at the vault; nightly-compile.sh cd's into KNOWLEDGE_REPO.
Environment=KNOWLEDGE_REPO=@VAULT_REPO@
ExecStart=@TOOLS_REPO@/scripts/nightly-compile.sh
# Give the run room; a large inbox + model latency can take a while.
TimeoutStartSec=1800

[Install]
WantedBy=default.target
