# scripts/knowledge-compile.path.in
# systemd USER path-unit template — generated by scripts/install.sh, which fills in
# @VAULT_REPO@ (the knowledge-vault repo). Edit this template, then re-run install.sh.
#
# Watches for a manual compile request the MCP server drops into inbox/ and starts
# knowledge-compile.service. Because the timer and this path unit start the SAME service,
# systemd runs only one instance at a time — that single-instance guarantee is the lock
# shared between the manual trigger and the nightly job.
[Unit]
Description=Watch for an on-demand knowledge-vault compile request

[Path]
# Fires when the sentinel exists; the service consumes (deletes) it, re-arming the watch.
PathExists=@VAULT_REPO@/inbox/.compile/request
Unit=knowledge-compile.service

[Install]
WantedBy=default.target
