# scripts/knowledge-compile@.timer.in
# systemd USER timer — generated PER VAULT by scripts/install.sh as knowledge-compile@<instance>.timer,
# which fills in @ONCALENDAR@ from KNOWLEDGE_COMPILE_ONCALENDAR (default: hourly) for that vault, so
# each vault can keep its own cadence. It triggers knowledge-compile@<instance>.service (same instance).
# The worker no-ops cheaply when the inbox is empty (no Claude run), so a tick only does real work when
# there are captures to compile. Change the cadence via that env var (or .env), then re-run install.sh.
[Unit]
Description=Scheduled %i vault inbox compile

[Timer]
OnCalendar=@ONCALENDAR@
# Run on next boot/login if the machine was off when a tick was due (collapses missed
# ticks into a single catch-up run).
Persistent=true
# Spread out so it doesn't fire at the exact same instant as other timers.
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
