# scripts/knowledge-resolve@.timer.in
# systemd USER timer — generated PER VAULT by scripts/install.sh as knowledge-resolve@<instance>.timer,
# which fills in @RESOLVE_ONCALENDAR@ from KNOWLEDGE_RESOLVE_ONCALENDAR (default: daily) for that vault.
# It triggers knowledge-resolve@<instance>.service (same instance). This is the light consumer pass; it
# short-circuits cheaply when there are no open issues, so a frequent cadence is fine. Change it via that
# env var (or .env), then re-run install.sh.
[Unit]
Description=Scheduled %i vault issue-resolve pass

[Timer]
OnCalendar=@RESOLVE_ONCALENDAR@
# Run on next boot/login if the machine was off when a tick was due.
Persistent=true
# Spread out so it doesn't fire at the same instant as other timers.
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
