# scripts/knowledge-compile@.service.in
# systemd USER service TEMPLATE — generated into ~/.config/systemd/user/knowledge-compile@.service
# by scripts/install.sh, which fills in @TOOLS_REPO@ (this repo, where the worker lives). One
# shared template serves every vault: the instance name %i (e.g. "default", "personal", "work")
# selects the per-vault config file %h/.config/knowledge-tools/%i.env, which carries this vault's
# KNOWLEDGE_REPO. Edit this template, then re-run install.sh.
#
# Runs one ephemeral inbox->wiki compile for instance %i. Started by that instance's timer and
# path units, which share this single per-instance service (so systemd serializes them — the
# cross-process lock for this vault):
#   - knowledge-compile@%i.timer → scheduled run (cadence via KNOWLEDGE_COMPILE_ONCALENDAR)
#   - knowledge-compile@%i.path  → on-demand run when the MCP drops inbox/.compile/request
#
[Unit]
Description=Compile the %i vault inbox into the wiki (ephemeral Claude run)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# This vault's config (KNOWLEDGE_REPO etc.). vault-compile.sh cd's into KNOWLEDGE_REPO; the lock
# keys off KNOWLEDGE_INSTANCE so different vaults run concurrently while this vault's jobs serialize.
EnvironmentFile=%h/.config/knowledge-tools/%i.env
Environment=KNOWLEDGE_INSTANCE=%i
ExecStart=@TOOLS_REPO@/scripts/vault-compile.sh
# Give the run room; a large inbox + model latency can take a while.
TimeoutStartSec=1800

[Install]
WantedBy=default.target
