# scripts/knowledge-synthesize@.timer.in
# systemd USER timer — generated PER VAULT by scripts/install.sh as knowledge-synthesize@<instance>.timer,
# which fills in @SYNTH_ONCALENDAR@ from KNOWLEDGE_SYNTHESIZE_ONCALENDAR (default: weekly) for that vault.
# It triggers knowledge-synthesize@<instance>.service (same instance). This is the heavy whole-corpus
# pass — keep it rare. Change the cadence via that env var (or .env), then re-run install.sh.
[Unit]
Description=Scheduled %i vault synthesis pass

[Timer]
OnCalendar=@SYNTH_ONCALENDAR@
# Run on next boot/login if the machine was off when a tick was due.
Persistent=true
# Spread out so it doesn't fire at the same instant as other timers.
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
