# scripts/knowledge-resolve@.service.in
# systemd USER service TEMPLATE — generated into ~/.config/systemd/user/knowledge-resolve@.service
# by scripts/install.sh, which fills in @TOOLS_REPO@ (this repo, where the worker lives). One
# shared template serves every vault: the instance name %i selects this vault's config file
# %h/.config/knowledge-tools/%i.env (KNOWLEDGE_REPO, and any per-vault KNOWLEDGE_GITHUB_REPO /
# KNOWLEDGE_REVIEW_CHANNEL). Edit this template, then re-run install.sh.
#
# Runs one ephemeral /resolve pass for instance %i: reads my answers on open judgment-call issues,
# applies the decisions to wiki/, and CLOSES the issues. Light + more frequent (cadence via
# KNOWLEDGE_RESOLVE_ONCALENDAR); often a no-op when nothing is answered. Like synthesize it runs
# `gh` from *inside* the Claude run, so PATH must reach gh and HOME must reach its auth.
[Unit]
Description=Resolve answered %i vault judgment-call issues into the wiki (closes them)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# This vault's config (KNOWLEDGE_REPO, per-vault KNOWLEDGE_GITHUB_REPO/REVIEW_CHANNEL). vault-job.sh
# cd's into KNOWLEDGE_REPO; the lock keys off KNOWLEDGE_INSTANCE so vaults don't block each other.
EnvironmentFile=%h/.config/knowledge-tools/%i.env
Environment=KNOWLEDGE_INSTANCE=%i
# gh (and claude) live in the nix/local profiles, which aren't on the default systemd user
# PATH. Claude shells out to `gh` for the issue calls, so those dirs must be here. HOME is set
# by the user manager and gives gh its auth at ~/.config/gh.
Environment=PATH=%h/.nix-profile/bin:%h/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# Optional machine-account identity. If this file exists (mode 0600, containing GH_TOKEN=<PAT>),
# gh authenticates as that bot account, so the issue comments it posts aren't authored by your
# personal account. The leading `-` makes it OPTIONAL: with no file, gh falls back to HOME's
# ~/.config/gh auth (commenting as you) — so this line is inert until you create the file. The
# token must never live in the repo; this keeps it out. See README "Bot account (optional)".
EnvironmentFile=-%h/.config/knowledge-tools/gh.env
ExecStart=@TOOLS_REPO@/scripts/vault-job.sh resolve
# Targeted opus pass over the named notes; lighter than synthesize.
TimeoutStartSec=1800

[Install]
WantedBy=default.target
